# teams/client-ops-team/ops/systemd/client-ops-commerce-handoff-watch.service
[Unit]
Description=Client Ops commerce handoff watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/ai-adam-agent-company/teams/client-ops-team
EnvironmentFile=-/opt/ai-adam-agent-company/teams/client-ops-team/.env
ExecStart=/opt/ai-adam-agent-company/teams/commerce-agent-team/.venv/bin/python scripts/watch_commerce_handoffs.py --once
