# teams/commerce-agent-team/ops/systemd/commerce-cleanup.timer
[Unit]
Description=Clean old commerce agent runs daily

[Timer]
OnCalendar=daily
Persistent=true
Unit=commerce-cleanup.service

[Install]
WantedBy=timers.target
