# teams/commerce-agent-team/ops/systemd/commerce-agents.service
[Unit]
Description=Commerce agent team run
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/opt/ai-adam-agent-company/teams/commerce-agent-team
EnvironmentFile=/opt/ai-adam-agent-company/teams/commerce-agent-team/.env
ExecStart=/opt/ai-adam-agent-company/teams/commerce-agent-team/.venv/bin/python scripts/run_agents.py --use-llm
